# config/autostart/Audacious.desktop
[Desktop Entry]
Exec=/usr/bin/audacious -p %U
Name=Audacious
OnlyShowIn=LXQt;
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
