# config/autostart/Cursor.desktop
[Desktop Entry]
Exec=xrdb ~/.Xdefaults
Name=Cursor
OnlyShowIn=LXQt;
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
