# config/autostart/Som.desktop
[Desktop Entry]
Exec=pactl set-sink-volume 0 60%
Name=Som
OnlyShowIn=LXQt;
Type=Application
Version=1.0
