# config/autostart/VLC.desktop
[Desktop Entry]
Exec=/usr/bin/vlc ~/Downloads/playlist.m3u8
Name=VLC
OnlyShowIn=LXQt;
Type=Application
Version=1.0
X-LXQt-Autostart-disabled=true
